# PI_ROS_WORKSPACES/test/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/home/pi/test/src/ros/rosbash/rosfish"
